# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ─── Servidor Node.js (Express) ──────────────────────────────────────────────
; Escucha en $PORT (asignado por Render). Si PORT no está definido, usa 3000.
[program:node]
command=node server.js
directory=/app
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ─── Servicio Python (FastAPI + Uvicorn) ─────────────────────────────────────
; Corre internamente en el puerto 8000.
; Node.js lo llama directamente via fetch('http://localhost:8000/...')
[program:python]
command=uvicorn stats:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5
