# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; Node.js Server (Express)
; Listens on $PORT (assigned by Render). Defaults to 3000 if not defined.
[program:node]
command=node server.js
directory=/app
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; Python Service (FastAPI + Uvicorn)
; Runs internally on port 8000.
; Accessed by Node.js via fetch('http://localhost:8000/...')
[program:python]
command=uvicorn stats:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5
